# BTWi-fi Autologin Openwrt Sevice/BTWi-fi_Autologin_Service
#!/bin/sh /etc/rc.common

##############- BTWi-fi Openwrt Autologin Service -##############
##############-- By: Aidan Macgregor (May 2022) --###############
# https://github.com/aidanmacgregor/BT_Wi-Fi_Autologin_MACRODROID-WISPr-HTTP_POST-HTTP_GET-OpenWRT
# (Tested On LEDE 17.01.7, OpenWrt 19.07.10 & 21.02.3)

############### INFO ###############

#### OpenWrt Install This Required Package:
# 	libustream-mbedtls
# (libustream-mbedtls Not Needed On OpenWRT 21.X as Wolf SSL included)

#### Copy This File To /etc/init.d/
# 	Use WinSCP To Tranfer File & Give 755 (Execute) permissions
# (If No Permissions You Get Error "Failed to execute "/etc/init.d/BTWi-fi_Autologin_Service start")

#### Manual Run & Stop (Stopping Also Signs The Account Out)
# 	Click Start/Stop on the service in LUCI (System > Startup)

#### Automatically Start On Boot
#	Enable the service in LUCI (System > Startup)

#### In The SETTINGS Section Choose Your Account Type & Add Email & Password
# (Reccomend Open This File With Notepad++)

#### Account Type:
# 	1 = BT Home Broadband
# 	2 = BT Buisness Broadband
# 	3 = BT Wi-Fi Account

###############- SETTINGS -###############

ACCOUNTTYPE=1
USERNAME=
PASSWORD=

###############- OPTIONAL -###############

PINGDNS=8.8.8.8
PING2URL=www.google.com
LOGPATH=/BTWi-fi_Autologin_Service_LOG.txt

##########################################
#####---- DO NOT EDIT BELOW HERE ----#####
##########################################

START=99
STOP=1
PIDFILE=/var/run/btwifi.pid

btwifi_loop(){
while true
do
if [[ "92" -ge "92" ]]; then
	echo "$(tail -92 $LOGPATH)" > $LOGPATH
fi
	if ! ping -c 1 -W 1 $PINGDNS 2>/dev/null >/dev/null	 
	then
		if ! ping -c 1 -W 1 $PING2URL 2>/dev/null >/dev/null			 
		then
			logger -t BTWi-fi_Autologin_Service "Ping 2 FAIL"
			echo "$(date) Ping 2 FAIL" >> $LOGPATH
			if [ "$ACCOUNTTYPE" = "1" ]
			then
				logger -t BTWi-fi_Autologin_Service "$(date) Offline, attempting login URL 1 (BT Home Broadband Account)"
				echo "$(date) Offline, attempting login URL 1 (BT Home Broadband Account)" >> $LOGPATH
				wget -T 2 -O /dev/null --post-data "username=$USERNAME&password=$PASSWORD" 'https://192.168.23.21:8443/tbbLogon'
			elif [ "$ACCOUNTTYPE" = "2" ]
			then
				logger -t BTWi-fi_Autologin_Service "$(date) Offline, attempting login URL 2 (BT Buisness Broadband Account)"
				echo "$(date) Offline, attempting login URL 2 (BT Buisness Broadband Account)" >> $LOGPATH
				wget -T 2 -O /dev/null --post-data "username=$USERNAME&password=$PASSWORD" 'https://192.168.23.21:8443/ante?partnerNetwork=btb'
			elif [ "$ACCOUNTTYPE" = "3" ]
			then
				logger -t BTWi-fi_Autologin_Service "$(date) Offline, attempting login URL 3 (BT Wi-Fi Account)"
				echo "$(date) Offline, attempting login URL 3 (BT Wi-Fi Account)" >> $LOGPATH
				wget -T 2 -O /dev/null --post-data "username=$USERNAME&password=$PASSWORD" 'https://192.168.23.21:8443/ante'
			fi
		fi
	fi
sleep 1
done
}

start() {
	logger -t BTWi-fi_Autologin_Service "$(date) BTWi-fi Autologin Service Started"
	echo "$(date) BTWi-fi Autologin Service Started" >> $LOGPATH
	[ -f $PIDFILE ] && [ ! -d /proc/`cat $PIDFILE` ] && rm $PIDFILE
	[ -f $PIDFILE ] && exit 1
	btwifi_loop &
	echo -n $! > $PIDFILE
}

stop() {
	kill `cat $PIDFILE`
	rm $PIDFILE
	wget -T 2 -O /dev/null 'https://192.168.23.21:8443/accountLogoff/home?confirmed=true'
	logger -t BTWi-fi_Autologin_Service "$(date) BTWi-fi Autologin Service Stopped Manually (Or Reboot)"
	echo "$(date) BTWi-fi Autologin Service Stopped Manually (Or Reboot)" >> $LOGPATH
}